# s-admin-watch.service
[Unit]
Description=ssh_shorty admin access-request watcher
Documentation=https://github.com/yadhusnair/ssh_shorty
After=network-online.target graphical-session.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/s-admin-watch
Restart=on-failure
RestartSec=15

# Make notify-send work from a systemd user service
Environment="DISPLAY=:0"
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus"

# Log to journald — view with: journalctl --user -u s-admin-watch -f
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
